# repos/core/initscripts/inittab
::sysinit:/etc/rc.boot
::respawn:/usr/bin/runsvdir /var/service
::shutdown:/etc/rc.shutdown
