# repos/core/initscripts/inittab
::sysinit:/etc/rc.boot
::respawn:/usr/bin/runsvdir /var/service
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/rc.shutdown
