# repos/extra/tlp/.files
drwxr-xr-x root/root    etc/
drwxr-xr-x root/root    etc/init.d/
-rwxr-xr-x root/root    etc/init.d/tlp.new
-rw-r--r-- root/root    etc/tlp.conf.new
drwxr-xr-x root/root    etc/tlp.d/
-rw-r--r-- root/root    etc/tlp.d/00-template.conf.new
-rw-r--r-- root/root    etc/tlp.d/README.new
drwxr-xr-x root/root    usr/
drwxr-xr-x root/root    usr/bin/
-rwxr-xr-x root/root    usr/bin/bluetooth
lrwxrwxrwx root/root    usr/bin/nfc -> bluetooth
-rwxr-xr-x root/root    usr/bin/run-on-ac
lrwxrwxrwx root/root    usr/bin/run-on-bat -> run-on-ac
-rwxr-xr-x root/root    usr/bin/tlp-rdw
-rwxr-xr-x root/root    usr/bin/tlp-stat
lrwxrwxrwx root/root    usr/bin/wifi -> bluetooth
lrwxrwxrwx root/root    usr/bin/wwan -> bluetooth
drwxr-xr-x root/root    usr/lib/
drwxr-xr-x root/root    usr/lib/NetworkManager/
drwxr-xr-x root/root    usr/lib/NetworkManager/dispatcher.d/
-rwxr-xr-x root/root    usr/lib/NetworkManager/dispatcher.d/99tlp-rdw-nm
drwxr-xr-x root/root    usr/sbin/
-rwxr-xr-x root/root    usr/sbin/tlp
drwxr-xr-x root/root    usr/share/
drwxr-xr-x root/root    usr/share/bash-completion/
drwxr-xr-x root/root    usr/share/bash-completion/completions/
lrwxrwxrwx root/root    usr/share/bash-completion/completions/bluetooth -> tlp
lrwxrwxrwx root/root    usr/share/bash-completion/completions/nfc -> tlp
lrwxrwxrwx root/root    usr/share/bash-completion/completions/run-on-ac -> tlp
lrwxrwxrwx root/root    usr/share/bash-completion/completions/run-on-bat -> tlp
-rw-r--r-- root/root    usr/share/bash-completion/completions/tlp
-rw-r--r-- root/root    usr/share/bash-completion/completions/tlp-rdw
lrwxrwxrwx root/root    usr/share/bash-completion/completions/tlp-stat -> tlp
lrwxrwxrwx root/root    usr/share/bash-completion/completions/wifi -> tlp
lrwxrwxrwx root/root    usr/share/bash-completion/completions/wwan -> tlp
drwxr-xr-x root/root    usr/share/fish/
drwxr-xr-x root/root    usr/share/fish/vendor_completions.d/
lrwxrwxrwx root/root    usr/share/fish/vendor_completions.d/bluetooth.fish -> tlp.fish
lrwxrwxrwx root/root    usr/share/fish/vendor_completions.d/nfc.fish -> tlp.fish
lrwxrwxrwx root/root    usr/share/fish/vendor_completions.d/run-on-ac.fish -> tlp.fish
lrwxrwxrwx root/root    usr/share/fish/vendor_completions.d/run-on-bat.fish -> tlp.fish
-rw-r--r-- root/root    usr/share/fish/vendor_completions.d/tlp-rdw.fish
-rw-r--r-- root/root    usr/share/fish/vendor_completions.d/tlp-stat.fish
-rw-r--r-- root/root    usr/share/fish/vendor_completions.d/tlp.fish
lrwxrwxrwx root/root    usr/share/fish/vendor_completions.d/wifi.fish -> tlp.fish
lrwxrwxrwx root/root    usr/share/fish/vendor_completions.d/wwan.fish -> tlp.fish
drwxr-xr-x root/root    usr/share/metainfo/
-rw-r--r-- root/root    usr/share/metainfo/de.linrunner.tlp.metainfo.xml
drwxr-xr-x root/root    usr/share/tlp/
drwxr-xr-x root/root    usr/share/tlp/bat.d/
-rwxr-xr-x root/root    usr/share/tlp/bat.d/05-thinkpad
-rwxr-xr-x root/root    usr/share/tlp/bat.d/10-thinkpad-legacy
-rwxr-xr-x root/root    usr/share/tlp/bat.d/15-lenovo
-rwxr-xr-x root/root    usr/share/tlp/bat.d/20-huawei
-rwxr-xr-x root/root    usr/share/tlp/bat.d/25-msi
-rwxr-xr-x root/root    usr/share/tlp/bat.d/30-samsung
-rwxr-xr-x root/root    usr/share/tlp/bat.d/35-lg
-rwxr-xr-x root/root    usr/share/tlp/bat.d/36-lg-legacy
-rwxr-xr-x root/root    usr/share/tlp/bat.d/40-sony
-rwxr-xr-x root/root    usr/share/tlp/bat.d/45-system76
-rwxr-xr-x root/root    usr/share/tlp/bat.d/50-toshiba
-rwxr-xr-x root/root    usr/share/tlp/bat.d/55-cros-ec
-rwxr-xr-x root/root    usr/share/tlp/bat.d/56-framework
-rwxr-xr-x root/root    usr/share/tlp/bat.d/60-macbook
-rwxr-xr-x root/root    usr/share/tlp/bat.d/65-dell
-rwxr-xr-x root/root    usr/share/tlp/bat.d/89-asus
-rwxr-xr-x root/root    usr/share/tlp/bat.d/90-generic
-rwxr-xr-x root/root    usr/share/tlp/bat.d/TEMPLATE
-rw-r--r-- root/root    usr/share/tlp/defaults.conf
-rw-r--r-- root/root    usr/share/tlp/deprecated.conf
drwxr-xr-x root/root    usr/share/tlp/func.d/
-rwxr-xr-x root/root    usr/share/tlp/func.d/05-tlp-func-pm
-rwxr-xr-x root/root    usr/share/tlp/func.d/10-tlp-func-cpu
-rwxr-xr-x root/root    usr/share/tlp/func.d/15-tlp-func-disk
-rwxr-xr-x root/root    usr/share/tlp/func.d/20-tlp-func-usb
-rwxr-xr-x root/root    usr/share/tlp/func.d/25-tlp-func-rf
-rwxr-xr-x root/root    usr/share/tlp/func.d/30-tlp-func-rf-sw
-rwxr-xr-x root/root    usr/share/tlp/func.d/35-tlp-func-batt
-rwxr-xr-x root/root    usr/share/tlp/func.d/40-tlp-func-bay
-rwxr-xr-x root/root    usr/share/tlp/func.d/45-tlp-func-gpu
-rwxr-xr-x root/root    usr/share/tlp/func.d/tlp-func-stat
-rw-r--r-- root/root    usr/share/tlp/rename.conf
-rwxr-xr-x root/root    usr/share/tlp/tlp-func-base
-rwxr-xr-x root/root    usr/share/tlp/tlp-pcilist
-rwxr-xr-x root/root    usr/share/tlp/tlp-readconfs
-rwxr-xr-x root/root    usr/share/tlp/tlp-usblist
drwxr-xr-x root/root    usr/share/zsh/
drwxr-xr-x root/root    usr/share/zsh/site-functions/
-rw-r--r-- root/root    usr/share/zsh/site-functions/_tlp
-rw-r--r-- root/root    usr/share/zsh/site-functions/_tlp-radio-device
-rw-r--r-- root/root    usr/share/zsh/site-functions/_tlp-rdw
-rw-r--r-- root/root    usr/share/zsh/site-functions/_tlp-run-on
-rw-r--r-- root/root    usr/share/zsh/site-functions/_tlp-stat
drwxr-xr-x root/root    var/
drwxr-xr-x root/root    var/lib/
drwxr-xr-x root/root    var/lib/tlp/
